# systemd/aes67-bidirectional.service
[Unit]
Description=AES67 NMOS Bidirectional Audio Node
After=network-online.target pipewire.service ptp4l.service
Wants=network-online.target
Requires=pipewire.service

[Service]
Type=simple
User=aes67
Group=audio
ExecStart=/usr/local/bin/rpi-aes67 -c /etc/rpi-aes67/config.json -m bidirectional
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log

# Real-time scheduling
LimitRTPRIO=95
LimitMEMLOCK=infinity

[Install]
WantedBy=multi-user.target
